# data/osec.service.in
[Unit]
Description=Forcefully sets passwords expired.
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@prefix@/bin/osec-wrapper

[Install]
WantedBy=multi-user.target
